# stow/autostart/autostart/zen.desktop
[Desktop Entry]
Exec=/run/media/randimu/General/Zen/zen/zen
Icon=
Name=zen
Path=
Terminal=False
Type=Application
